# config/templates/startup_script.sh
#!/bin/bash
# chkconfig: 2345 20 80
# description: Description comes here....

# Source function library.
